# etc/moodns.service
[Unit]
Description=Minimal multicast DNS server
Documentation=man:moodns(1)

[Service]
Type=simple
ExecStart=/usr/bin/moodns
